# system/systemd/beo-source-news.service
[Unit]
Description=BeoSound5c News Source (Guardian API article browser)
After=network-online.target beo-input.service beo-router.service
Wants=network-online.target
OnFailure=beo-notify-failure@%n.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=bs5
Group=bs5
WorkingDirectory=/home/bs5/beosound5c/services
EnvironmentFile=-/etc/beosound5c/secrets.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/bs5/beosound5c/services/sources/news.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
