# system/systemd/beo-player-local.service
[Unit]
Description=BeoSound5c Local Player (mpv, playback commands)
After=network.target
Wants=network.target
OnFailure=beo-notify-failure@%n.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=notify
User=bs5
Group=bs5
WorkingDirectory=/home/bs5/beosound5c/services
EnvironmentFile=-/etc/beosound5c/secrets.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/bs5/beosound5c/services/players/local.py
Restart=on-failure
RestartSec=5
WatchdogSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
